# prov/cxi/src/cix_driver/50-cxi-driver.rules
KERNEL=="cxi_ss1", SUBSYSTEM=="module", ACTION=="add", RUN+="/sbin/modprobe --use-blacklist -a cxi_user cxi_eth"
